# deployment/ansible/templates/meeting-backend.service.j2
[Unit]
Description=Meeting Hotspot Backend
After=network.target

[Service]
WorkingDirectory={{ backend_repo_dir }}
EnvironmentFile={{ backend_repo_dir }}/.env
ExecStart={{ backend_repo_dir }}/.venv/bin/litestar --app app.main:create_app run --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
User=pi
Group=pi

[Install]
WantedBy=multi-user.target

